# deploy/interview-web.service
[Unit]
Description=Interview Bot — Web Admin (uvicorn)
Documentation=https://github.com/<your-repo>
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/srv/interview

# Читает переменные из .env; файл должен иметь права 600 (только владелец).
EnvironmentFile=/srv/interview/.env

ExecStart=/srv/interview/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1

# Перезапускать автоматически при сбое; ждать 5 с перед повтором.
Restart=always
RestartSec=5

# Ограничения безопасности: процесс не может писать в /usr, /boot и т.д.
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
